# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/log/supervisor/supervisord.pid
user=root

[program:engine]
command=python3 /app/src/engine_server.py 5001
priority=10
startsecs=2
autorestart=true
stdout_logfile=/var/log/supervisor/engine.log
stderr_logfile=/var/log/supervisor/engine.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:server]
command=bun src/server.ts
priority=20
startsecs=3
autorestart=true
environment=ENGINE_URL=http://127.0.0.1:5001,PORT=3280
stdout_logfile=/var/log/supervisor/server.log
stderr_logfile=/var/log/supervisor/server.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
